# tools/boot/rootfs/etc/profile.d/zzz_fancontrol.sh
# Display on connection

/bin/echo "======================================================="
/bin/echo "=                                                     ="
/bin/echo "=                  Fancontrol Linux                   ="
/bin/echo "=                                                     ="
/bin/echo "======================================================="

alias ll='ls -la'

